# packaging/systemd/prune-guard.timer
[Unit]
Description=Prune Guard Cleanup Daemon Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true
Unit=prune-guard.service

[Install]
WantedBy=timers.target
